# units-n-schema/snek-schedule.timer
[Unit]
Description=Run snek scheduler every minute

[Timer]
OnCalendar=*:0/1

