# units-n-schema/snek-schedule.service
[Unit]
Description=Run SnekServer scheduler

[Service]
Type=oneshot
ExecStart=/usr/bin/php /home/snek/snek/artisan schedule:run
